# bin/services/grafana.service
[Unit]
Description=Grafana
Wants=network-online.target
After=network-online.target
 
[Service]
User=pi
Group=pi
Type=simple
WorkingDirectory=/opt/grafana
RuntimeDirectory=grafana
ExecStart=/opt/grafana/bin/grafana-server
 
[Install]
WantedBy=multi-user.target

